# supervisor_conf.d/spotlite-syncer-worker.conf
[program:spotlite-syncer-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/vagrant/Code/spotlite/artisan queue:work --tries=3 --queue=syncing
autostart=true
autorestart=true
user=vagrant
numprocs=2
startsecs=0
redirect_stderr=true
stdout_logfile=/home/vagrant/Code/spotlite/storage/logs/spotlite-syncer-worker.log
